# ARCHITECTURAL/ArchiGPT_0.2/datasets/initial/4_Green_Thumb/green_thumb_backend-main/startup.service
[Unit]
Description=startup
After=network.target

[Service]
User=bob
WorkingDirectory=/home/bob/green_thumb_backend
#Environment="ENV_VAR_NEEDED=value"
#ExecStartPre=/bin/bash -c `/usr/bin/env`
ExecStart=/usr/bin/docker-compose up --build

[Install]
WantedBy=multi-user.target